# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for SITENAME"

start on net-device-up
stop on shutdown

respawn

setuid metatomato
chdir /home/metatomato/SITES/SITENAME/source
exec ../virtualenv/bin/gunicorn --bind unix:/tmp/SITENAME.socket iGLouWeb.wsgi:application
